# contrib/meridian.service
[Unit]
Description=Meridian CRDT sync server
Documentation=https://github.com/Chahine-tech/meridian
After=network.target
Wants=network.target

[Service]
Type=simple
User=meridian
Group=meridian
WorkingDirectory=/var/lib/meridian

# Environment file — copy from /etc/meridian/env.example
EnvironmentFile=/etc/meridian/env

ExecStart=/usr/local/bin/meridian

Restart=on-failure
RestartSec=5s

# Allow a large number of open file descriptors (WebSocket connections)
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/meridian

[Install]
WantedBy=multi-user.target
